# systemd units/daemon0.socket
[Unit]
Description=Network Yeller TCP Activating Socket

[Socket]
ListenStream=8027
Accept=yes

[Install]
WantedBy=sockets.target